# zookeeper
#/bin/sh  
### BEGIN INIT INFO
# Provides:          zookeeper
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: centralized coordination service
# Description:       ZooKeeper is a centralized service for maintaining
#                    configuration information, naming, providing distributed
#                    synchronization, and providing group services.
### END INIT INFO

case $1 in
          start)  /opt/zookeeper-3.4.7/bin/zkServer.sh start   ;;
          stop)  /opt/zookeeper-3.4.7/bin/zkServer.sh stop;;
          status)  /opt/zookeeper-3.4.7/bin/zkServer.sh status;;
          restart)  /opt/zookeeper-3.4.7/bin/zkServer.sh restart;;
              *)  echo "require start|stop|status|restart"  ;;
esac              